# dot_config/systemd/user/hyprland-autostart-whatsapp-web.service
[Unit]
Description=Hyprland autostart: WhatsApp Web
PartOf=graphical-session.target
Wants=caelestia-shell.service
After=caelestia-shell.service
# Forward into the already-running Trivalent instance instead of racing it for
# the profile. trivalent.service only reports ready once its singleton is live.
Wants=hyprland-autostart-trivalent.service
After=hyprland-autostart-trivalent.service

[Service]
Type=exec
ExecStart=/usr/bin/uwsm app -- /usr/lib64/trivalent/trivalent.sh --profile-directory=Default --app-id=hnpfjngllnobngcgfapefoaidbinmjnm

[Install]
WantedBy=graphical-session.target
